# docker/config/cron.d/crowdnews
5 18 * * 0 symfony /app/bin/console app:newsletter:prepare
0 6 * * 1 symfony /app/bin/console app:newsletter:send
